# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:rstudio]
command=/usr/lib/rstudio-server/bin/rserver --server-daemonize=0 --www-port=8787
autostart=true
autorestart=true
stdout_logfile=/var/log/rstudio.log
stderr_logfile=/var/log/rstudio.log

[program:wetty]
command=wetty --port 7681 --base /terminal --ssh-host localhost --ssh-port 22 --ssh-user rstudio
autostart=true
autorestart=true
stdout_logfile=/var/log/wetty.log
stderr_logfile=/var/log/wetty.log

[program:filebrowser]
command=filebrowser --database /etc/filebrowser.db
autostart=true
autorestart=true
stdout_logfile=/var/log/filebrowser.log
stderr_logfile=/var/log/filebrowser.log

[program:update-api]
command=python3 /usr/local/bin/update_api.py
autostart=true
autorestart=true
stdout_logfile=/var/log/update-api.log
stderr_logfile=/var/log/update-api.log