# s6/update-api/run
#!/bin/bash
exec python3 /usr/local/bin/update_api.py